# gnu+X+linux/systemd/apodesktop.service
[Unit]
Description=Change wallpapers to NASA Astronomy Picture of the Day

[Service]
ExecStart=/usr/local/bin/apodesktop
